# yocto/libreswitch/meta-distro-libreswitch/recipes-asic/opennsl/opennsl-odp/bcm.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="bdev", RUN+="/etc/udev/rules.d/bcm_devices.sh"
