# yocto/libreswitch/meta-distro-libreswitch/recipes-ops/utils/ops-broadview/ops-broadview.service

[Unit]
Description=LibreSwitch BroadView Daemon
After=ovsdb-server.service ops-init.service

[Service]
Type=simple
ExecStartPre=/bin/rm -f /var/run/openvswitch/ops-broadview.pid
ExecStart= /usr/sbin/ops-broadview --pidfile
Restart=on-abort

[Install]
WantedBy=multi-user.target
