# yocto/openhalon/meta-platform-openhalon-as5712/recipes-core/systemd/systemd/onie-install.target
[Unit]
Description=ONIE Installer mode
Documentation=man:systemd.special(7)
DefaultDependencies=no
Requires=onie-install.service
After=onie-install.service
AllowIsolate=yes

[Install]
Alias=onie-install.target
