# yocto/halon/meta-hp/recipes-devtools/t32server/files/t32.service
[Unit]
Description=T32 Debug Server
ConditionPathExists=|/usr/bin/t32server

[Service]
ExecStart=/usr/bin/t32server

[Install]
WantedBy=getty.target
