# yocto/openswitch/meta-foss-openswitch/recipes-connectivity/openssh/openssh/sshd@.service
[Unit]
Description=OpenSSH Per-Connection Daemon
Wants=sshdgenkeys.service
After=sshdgenkeys.service

[Service]
Environment=LD_PRELOAD=/usr/lib/cli/plugins/libaaa_utils_lib.so
ExecStart=-@SBINDIR@/sshd -i
ExecReload=@BASE_BINDIR@/kill -HUP $MAINPID
StandardInput=socket
StandardError=syslog
KillMode=process
