# yocto/openswitch/meta-distro-openswitch/recipes-ops/utils/ops-kdump/kdump.target
[Unit]
Description=KDump Recovery Target

Requires=sysinit.target
After=sysinit.target
Wants=sockets.target timers.target paths.target slices.target
After=sockets.target paths.target slices.target
