# yocto/openswitch/meta-distro-openswitch/recipes-ops/l2/ops-lacpd/ops-lacpd.service
# Copyright (C) 2015 Hewlett Packard Enterprise Development LP

[Unit]
Description=OpenSwitch LACP Daemon
After=ovsdb-server.service ops-init.service

[Service]
Type=forking
CPUSchedulingPolicy=rr
CPUSchedulingPriority=3
ExecStartPre=/bin/rm -f /var/run/openvswitch/ops-lacpd.pid
ExecStart=/sbin/ip netns exec swns /usr/bin/ops-lacpd --detach --pidfile -vSYSLOG:INFO
PIDFile=/var/run/openvswitch/ops-lacpd.pid
Restart=on-abort

[Install]
WantedBy=multi-user.target
