# yocto/openswitch/meta-distro-openswitch/recipes-ops/platform/ops-reboot/switch_reboot@.service
# Copyright (C) 2016 Hewlett Packard Enterprise Development LP

[Unit]
Description=OpenSwitch system reboot script

[Service]
Type=oneshot
ExecStart=/usr/sbin/ops-reboot %i
